# containers/update/crontab
* * * * * python /update_data.py
0 2 */15 * * python /update_models.py
